# scripts/resource-monitor.service
[Unit]
Description=Quant Engine Resource Monitor (VPS RAM/CPU guardian)
After=default.target
Documentation=file:///home/willy/invest/quant-engine/.claude/skills/process-cleanup/SKILL.md

[Service]
Type=simple
WorkingDirectory=/home/willy/invest/quant-engine
ExecStart=/usr/bin/python3 /home/willy/invest/quant-engine/scripts/resource_monitor.py
Restart=always
RestartSec=30
# Low priority so monitor never competes with trading system
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=5
# Guard against runaway memory use in the monitor itself
MemoryMax=128M
# Logs go to journalctl --user-unit=quant-resource-monitor AND the rotating file
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quant-resource-monitor

[Install]
WantedBy=default.target
